# scripts/fin1-seneschal-tunnel.service
# Persistent reverse tunnel fin1 → fin4 exposing local chain RPCs to the
# data-api host: monerod (18081), zebra (8232) and NFPT (3555) appear on
# fin4's 127.0.0.1. The data-api, its pollers and the ops-monitor probes
# all depend on these forwards.
#
# Installed on fin1 as /etc/systemd/system/fin1-seneschal-tunnel.service:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now fin1-seneschal-tunnel.service
#
# Replaces the per-user unit (~/.config/systemd/user/seneschal-chain-tunnel
# .service, retired 2026-06-12): with Linger=no a user unit only runs while
# rotwang has a session, so an unattended reboot left the tunnel down.
# A system unit starts at boot unconditionally.
#
# Runs as rotwang with the dedicated tunnel key — restricted on fin4 with
# `restrict,port-forwarding,command="echo seneschal-tunnel-only"`, so a
# leaked key can only open this same tunnel: no shell, no exec.
#
# ExitOnForwardFailure + Restart=always: a half-up tunnel (binds taken /
# connection dropped) dies and retries rather than lingering. fin4's sshd
# runs ClientAliveInterval 30 / CountMax 3 (sshd_config.d/99-clientalive
# .conf) so a dead client's stale listeners are reaped within ~90s —
# without it a SIGKILL'd client wedges the rebind forever.
# StartLimitIntervalSec=0: never stop retrying — after a long fin4 outage
# the tunnel must come back on its own.

[Unit]
Description=Seneschal reverse tunnel to fin4 (monerod/zebra/NFPT forwards)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=exec
User=rotwang
ExecStart=/usr/bin/ssh -NT \
	-i /home/rotwang/.ssh/seneschal_tunnel_ed25519 \
	-o IdentitiesOnly=yes \
	-o BatchMode=yes \
	-o ExitOnForwardFailure=yes \
	-o ServerAliveInterval=30 \
	-o ServerAliveCountMax=3 \
	-o ConnectTimeout=15 \
	-o StrictHostKeyChecking=accept-new \
	-o StreamLocalBindUnlink=yes \
	-R 127.0.0.1:18081:127.0.0.1:18081 \
	-R 127.0.0.1:8232:127.0.0.1:8232 \
	-R 127.0.0.1:3555:127.0.0.1:3555 \
	root@fin4
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
